# bin/git-version-update.service
[Unit]
Description=RemoteLabz Git Version Update Service
After=network.target

[Service]
Type=oneshot
ExecStart=/bin/bash /opt/remotelabz/bin/git-version-update.sh
User=www-data
Group=www-data
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target