# bin/remotelabz-git-version-update.timer
[Unit]
Description=RemoteLabz Git Version Update Timer
Requires=remotelabz-git-version-update.service

[Timer]
# Exécuter au démarrage
OnBootSec=30s
# Puis toutes les heures
OnUnitActiveSec=1h
# Persist le timer à travers les redémarrages
Persistent=true

[Install]
WantedBy=timers.target