# bin/remotelabz.service
[Unit]
Description=Remotelabz Messages Consumer Service
After=network.target

[Service]
Type=simple
Restart=always
#Restart=on-failure
RestartSec=5
ExecStartPre=/usr/bin/env php /opt/remotelabz/bin/remotelabz-network-up
ExecStart=/usr/bin/env php /opt/remotelabz/bin/console messenger:consume front --memory-limit=128M
ExecStop=/usr/bin/env php /opt/remotelabz/bin/console messenger:stop-workers
MemoryAccounting=true
MemoryHigh=128M
MemoryMax=256M

[Install]
WantedBy=multi-user.target